# packaging/wavemux.service
[Unit]
Description=WaveMux Audio Mixer Daemon
After=pipewire.service pipewire-pulse.service
Wants=pipewire-pulse.service

[Service]
Type=simple
ExecStart=/usr/bin/wavemuxd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
